# etc/udev/99-kill-marimo.rules
ACTION=="poweroff", RUN+="/usr/bin/pkill -f marimo"
ACTION=="reboot", RUN+="/usr/bin/pkill -f marimo"

